# scripts/microgrid-dashboard.service
[Unit]
Description=Microgrid Streamlit Dashboard
After=network.target

[Service]
User=amplifiedengr
WorkingDirectory=/opt/microgrid
ExecStart=/usr/local/bin/streamlit run dashboard/dashboard_app.py \
    --server.port=8501 \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --server.enableCORS=false \
    --server.enableXsrfProtection=false
Restart=always
RestartSec=10
MemoryMax=1500M
Environment=GCP_PROJECT_ID=microgrid-demo

[Install]
WantedBy=multi-user.target
